# celery/flower.service
[Unit]
Description=Flower celery service

[Service]
ExecStart=/bin/bash -c "/EXCITE/celery/flower.sh"

[Install]
WantedBy=multi-user.target